# bin/his/mac_cron/self.investment.data.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>  
    <key>Label</key>  
    <string>self.investment.data</string>  
    <key>ProgramArguments</key>  
      <array>  
        <string>/Users/kay/SourceCodes/stock_ays/bin/daily/investment.sh</string>  
      </array>  
    <key>StartCalendarInterval</key>  
      <dict>  
            <key>Minute</key>  
            <integer>25</integer>  
            <key>Hour</key>  
            <integer>15</integer>  
      </dict>  
    <key>StandardOutPath</key>  
      <string>/dev/null</string>  
    <key>StandardErrorPath</key>  
      <string>/dev/null</string>  
  </dict>  
</plist>  
